# etc/crontab
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
*/3 * * * * cd $HOME/projects/bbbike/ && world/bin/extract.pl
5 * * * * cd $HOME/projects/bbbike && make -s extract-gc7 extract-gc
MAILTO=wosch
# m h  dom mon dow   command
15,45 * * * 4     cd $HOME/projects/bbbike/ && world/bin/planet-download-alert
15 1,3,7,10 * * 5 cd $HOME/projects/bbbike/ && world/bin/planet-download-alert

